# scripts/net/keep_internet.service
[Unit]
After=network.service

[Service]
ExecStart=/home/scripts/net/keep_internet.sh

[Install]
WantedBy=default.target
