# scripts/net/keep_internet.service
[Unit]
After=network.service

[Service]
ExecStart=$HAIOT_DIR/net/keep_internet.sh

[Install]
WantedBy=default.target
