# scripts/osinstall/ubuntu/etc/systemd/system/activate-audio-amp.service
[Unit]
After=mpd@.service

[Service]
ExecStart=/home/scripts/audio/activate-audio-amp.sh loop

[Install]
WantedBy=default.target
