# scripts/osinstall/ubuntu/etc/systemd/system/mpd@beci.socket
[Socket]
ListenStream=6602

[Install]
WantedBy=sockets.target
