# scripts/osinstall/ubuntu/etc/systemd/system/mpd@baie.socket
[Socket]
ListenStream=6604

[Install]
WantedBy=sockets.target
