# scripts/osinstall/ubuntu/etc/systemd/system/mpd@living.socket
[Socket]
ListenStream=6600

[Install]
WantedBy=sockets.target
