# scripts/osinstall/ubuntu/etc/systemd/system/mpd@dormitor.socket
[Socket]
ListenStream=6603

[Install]
WantedBy=sockets.target
