# scripts/osinstall/ubuntu/etc/systemd/system/record-audio.service
[Unit]
After=mpd@.service

[Service]
ExecStart=/home/scripts/audio/record-audio.sh

[Install]
WantedBy=default.target
