# scripts/osinstall/ubuntu/etc/systemd/system/thd.service
[Unit]
After=mpd@.service

[Service]
ExecStart=/home/scripts/input/thd-daemon.sh

[Install]
WantedBy=default.target
